# mkosi.profiles/runtime/mkosi.extra/usr/lib/systemd/system/containerd.service.d/10-katl-runtime.conf
[Unit]
Requires=var.mount
After=var.mount
Before=katl-kubeadm-ready.target
RequiresMountsFor=/var/lib/containerd
