# mkosi.profiles/runtime/mkosi.extra/usr/lib/systemd/system/katl-boot-complete.target
[Unit]
Description=Katl boot-complete promotion point
Documentation=man:systemd.target(5)
Requires=katl-boot-health.service
After=katl-boot-health.service

[Install]
WantedBy=multi-user.target
