# mkosi.profiles/runtime/mkosi.extra/usr/lib/systemd/system/katl-boot-health.service
[Unit]
Description=Record successful Katl boot health
Documentation=man:systemd.service(5)
Requires=katl-runtime-handoff-status.service katlc-agent.service systemd-networkd.service
Wants=sshd.service
After=katl-runtime-handoff-status.service katlc-agent.service systemd-networkd.service sshd.service
Before=katl-boot-complete.target
RequiresMountsFor=/var/lib/katl

[Service]
Type=oneshot
StandardOutput=journal+console
SyslogIdentifier=katl-boot-health
ExecStart=/usr/lib/katl/runtime/katl-boot-health --root=/ --result=success --reason=katl-boot-complete.target

[Install]
RequiredBy=katl-boot-complete.target
