# mkosi.profiles/runtime/mkosi.extra/usr/lib/systemd/system/katl-boot-deadman.service
[Unit]
Description=Fail Katl boot health after deadline
Documentation=man:systemd.service(5)
Requires=var.mount
After=var.mount
RequiresMountsFor=/var/lib/katl

[Service]
Type=oneshot
StandardOutput=journal+console
SyslogIdentifier=katl-boot-deadman
ExecStart=/usr/lib/katl/runtime/katl-boot-health --root=/ --result=timeout --reason=katl-boot-health-deadline-expired --request-reboot
