# mkosi.profiles/runtime/mkosi.extra/usr/lib/systemd/system/katl-runtime-handoff-status.service
[Unit]
Description=Record Katl runtime handoff status
Documentation=man:systemd.service(5)
Requires=katlc-agent.service
After=katlc-agent.service
Before=katl-boot-complete.target
RequiresMountsFor=/var/lib/katl

[Service]
Type=oneshot
StandardOutput=journal+console
SyslogIdentifier=katl-runtime-status
ExecStart=/usr/lib/katl/runtime/katl-runtime-status --root=/

[Install]
RequiredBy=katl-boot-complete.target
