# mkosi.profiles/runtime/mkosi.extra/usr/lib/systemd/system/katl-boot-deadman.timer
[Unit]
Description=Katl boot health deadline
Documentation=man:systemd.timer(5)

[Timer]
OnBootSec=10min
Unit=katl-boot-deadman.service

[Install]
WantedBy=timers.target
