# mkosi.profiles/installer-image/mkosi.extra/etc/systemd/system/sshd.service.d/10-katl-disabled.conf
[Unit]
ConditionPathExists=/etc/katl/installer-ssh.enabled
