# mkosi.profiles/runtime/mkosi.extra/usr/lib/systemd/system/kubelet.service.d/10-katl-runtime.conf
[Unit]
ConditionPathExists=/var/lib/kubelet/config.yaml
Requires=containerd.service etc-kubernetes.mount
After=var.mount containerd.service etc-kubernetes.mount
Before=katl-kubeadm-ready.target
RequiresMountsFor=/var/lib/kubelet /etc/kubernetes

[Service]
Environment=KUBELET_EXTRA_ARGS=
